# runtime/lid-monitor.service.in
[Unit]
Description=Hyprland Lid Switch Monitor
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=HYPRLAND_INSTANCE_SIGNATURE
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=exec
ExecStartPre=$HYPRCTL_BIN monitors
ExecStart=$SCRIPTS_DIR/lid-monitor.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arch-lidswitch
Restart=on-failure
RestartPreventExitStatus=1
RestartSec=2

[Install]
WantedBy=graphical-session.target
